# services/dnscrypt-proxy/run
#!/usr/bin/with-contenv bash

s6-echo "Starting dnscrypt-proxy"
/opt/dnscrypt-proxy/dnscrypt-proxy
